# utils/systemd/inbacoindd.service
[Unit]
Description=InbaCoin Full Node
After=network.target

[Service]
User=inbacoin
Group=inbacoin
WorkingDirectory=~
RuntimeDirectory=inbacoin

# Clearnet config
#
Type=forking
PIDFile=/run/inbacoin/inbacoind.pid
ExecStart=/usr/bin/inbacoind --config-file /etc/inbacoind.conf \
    --detach --pidfile /run/inbacoin/inbacoind.pid

# Tor config
#
## We have to use simple, not forking, because we cannot pass --detach
## because stderr/stdout is not available when detached, but torsocks
## attempts to write to it, and fails with 'invalid argument', causing
## monerod to fail.
#Type=simple
#Environment=DNS_PUBLIC=tcp
## The following is needed only when accessing wallet from a different
## host in the LAN, VPN, etc, the RPC must bind to 0.0.0.0, but
## by default torsocks only allows binding to localhost.
#Environment=TORSOCKS_ALLOW_INBOUND=1
#ExecStart=/usr/bin/torsocks /usr/bin/inbacoind --config-file /etc/inbacoind.conf \
#    --non-interactive

Restart=always

[Install]
WantedBy=multi-user.target
